# scripts/etc/cron.d/vm-sandboxes
# Copyright (c) Joseph Wolff and the Libre Group 2008
#
# /etc/cron.d/vm-sandboxes: restart sandbox VMs even so often with  vm -x

#  This loops through the active VMs and does a restart

*/15 * * * *     www-data   /usr/local/bin/vm -x
